# rootfs/etc/crontabs/root
# do daily/weekly/monthly maintenance
# min	hour	day	month	weekday	command

*/15	*	*	*	*	/srv/sbin/run-cron 15min
0	*	*	*	*	/srv/sbin/run-cron hourly
0	2	*	*	*	/srv/sbin/run-cron daily
0	3	*	*	6	/srv/sbin/run-cron weekly
0	5	1	*	*	/srv/sbin/run-cron monthly
