# rootfs/etc/inittab
# /etc/inittab init(8) configuration for BusyBox
# https://github.com/mirror/busybox/blob/master/examples/inittab
#

::sysinit:/usr/bin/wkit-init
